# packaging/postinst
#!/usr/bin/env bash
# Post-installation script for kentik-device-onboarder.
# Compatible with both DEB (argument: "configure") and RPM (argument: "1" install / "2" upgrade).

set -e

SERVICE_USER=kentik-onboarder
SERVICE_GROUP=kentik-onboarder
INSTALL_DIR=/opt/kentik-device-onboarder
CONFIG_DIR=/etc/kentik-device-onboarder
STATE_DIR=/var/lib/kentik-device-onboarder
SERVICE_NAME=kentik-device-onboarder.service

discover_kproxy_credentials() {
    local pid environ_path line key value
    local found_email=""
    local found_token=""

    pid=$(pgrep -n kproxy 2>/dev/null || true)
    if [ -z "$pid" ]; then
        return 1
    fi

    environ_path="/proc/$pid/environ"
    if [ ! -r "$environ_path" ]; then
        return 1
    fi

    while IFS='=' read -r key value; do
        case "$key" in
            KENTIK_API_EMAIL) found_email="$value" ;;
            KENTIK_API_TOKEN) found_token="$value" ;;
        esac
    done < <(tr '\0' '\n' < "$environ_path")

    if [ -n "$found_email" ] && [ -n "$found_token" ]; then
        KPROXY_KENTIK_API_EMAIL="$found_email"
        KPROXY_KENTIK_API_TOKEN="$found_token"
        return 0
    fi

    return 1
}

escape_sed_replacement() {
    printf '%s' "$1" | sed -e 's/[\\&|]/\\&/g'
}

populate_credentials_in_config() {
    local config_file="$1"
    local email_escaped token_escaped

    if ! discover_kproxy_credentials; then
        echo "kproxy credentials not found in process environment; leaving API credentials unchanged"
        return 0
    fi

    email_escaped=$(escape_sed_replacement "$KPROXY_KENTIK_API_EMAIL")
    token_escaped=$(escape_sed_replacement "$KPROXY_KENTIK_API_TOKEN")

    sed -i \
        -e "s|^KENTIK_API_EMAIL=.*$|KENTIK_API_EMAIL=${email_escaped}|" \
        -e "s|^KENTIK_API_TOKEN=.*$|KENTIK_API_TOKEN=${token_escaped}|" \
        "$config_file"

    echo "populated KENTIK_API_EMAIL and KENTIK_API_TOKEN from kproxy environment"
}

_configure() {
    if ! getent group "$SERVICE_GROUP" >/dev/null 2>&1; then
        groupadd --system "$SERVICE_GROUP"
    fi

    if ! id -u "$SERVICE_USER" >/dev/null 2>&1; then
        useradd \
            --system \
            --gid "$SERVICE_GROUP" \
            --home-dir "$INSTALL_DIR" \
            --no-create-home \
            --shell /usr/sbin/nologin \
            "$SERVICE_USER"
    fi

    install -d -m 0750 -o "$SERVICE_USER" -g "$SERVICE_GROUP" "$STATE_DIR"

    chown root:"$SERVICE_GROUP" "$CONFIG_DIR"
    chmod 0750 "$CONFIG_DIR"

    if [ ! -f "$CONFIG_DIR/onboarder.env" ]; then
        install -m 0640 -o "$SERVICE_USER" -g "$SERVICE_GROUP" \
            "$CONFIG_DIR/onboarder.env.example" \
            "$CONFIG_DIR/onboarder.env"
        populate_credentials_in_config "$CONFIG_DIR/onboarder.env"
        echo "Created initial configuration file: $CONFIG_DIR/onboarder.env"
        echo "Edit this file and add your Kentik credentials before starting the service."
    fi

    chown root:"$SERVICE_GROUP" "$INSTALL_DIR"
    chown "$SERVICE_USER":"$SERVICE_GROUP" "$INSTALL_DIR/kentik_device_onboarder.py"

    if command -v systemctl >/dev/null 2>&1; then
        systemctl daemon-reload
        systemctl enable "$SERVICE_NAME"
        echo "Service enabled. Start with: systemctl start $SERVICE_NAME"
    fi
}

# DEB passes "configure"; RPM passes "1" (install) or "2" (upgrade).
case "${1:-}" in
    configure|1|2)
        _configure
        ;;
    abort-upgrade|abort-remove|abort-deconfigure)
        ;;
    *)
        # Unknown argument — still run configuration so bare invocation works.
        _configure
        ;;
esac

exit 0
